# osx/Library/LaunchAgents/com.ibm.de.suhlig.brew-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ibm.de.suhlig.brew-cleanup</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/brew</string>
      <string>cleanup</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>00</integer>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/suhlig/com.ibm.de.suhlig.brew-cleanup.err</string>
    <key>StandardOutPath</key>
    <string>/Users/suhlig/com.ibm.de.suhlig.brew-cleanup.out</string>
  </dict>
</plist>
